# apps/public/tools/sysstat-9.0.4/sysstat.cron.hourly.in
#!/bin/sh
# Run system activity accounting tool every @CRON_INTERVAL@ minutes
@SA_LIB_DIR@/sa1 -S DISK @CRON_INTERVAL_SEC@ @CRON_COUNT@ &
