# apps/public/tools/sysstat-9.0.4/sysstat.cron.hourly
#!/bin/sh
# Run system activity accounting tool every 10 minutes
/usr/lib/sa/sa1 -S DISK 600 6 &
